# packer/projects/aws-marketplace-ami/files/ade.service
[Unit]
Description=ADE - Agent Development Environment
After=network.target letta.service
Requires=letta.service
Wants=network.target

[Service]
Type=simple
User=ade
Group=ade
WorkingDirectory=/opt/ade
ExecStart=/opt/ade/run_ade.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ade

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ade /var/log/ade /tmp

# Environment
Environment=ADE_CONFIG_PATH=/etc/ade/config.yaml

[Install]
WantedBy=multi-user.target
