# packer/projects/aws-marketplace-ami/files/letta-platform.service
[Unit]
Description=Letta Platform - Containerized AI Agent Platform
Requires=docker.service
After=docker.service network.target
Wants=network.target

[Service]
Type=forking
RemainAfterExit=yes
WorkingDirectory=/opt/letta-platform
Environment=COMPOSE_PROJECT_NAME=letta-platform

# Start command
ExecStart=/usr/bin/docker-compose --env-file=.env up -d

# Stop command
ExecStop=/usr/bin/docker-compose --env-file=.env down

# Reload command
ExecReload=/usr/bin/docker-compose --env-file=.env restart

# Health check
ExecStartPost=/bin/bash -c 'sleep 30 && /opt/health_check.sh'

# Restart policy
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=120

# Security settings
User=root
Group=root

# Standard output/error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=letta-platform

[Install]
WantedBy=multi-user.target
