# packer/projects/aws-marketplace-ami/files/letta.service
[Unit]
Description=Letta AI Agent Platform
After=network.target postgresql.service
Requires=postgresql.service
Wants=network.target

[Service]
Type=simple
User=letta
Group=letta
WorkingDirectory=/opt/letta
ExecStart=/opt/letta/run_letta.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=letta

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/letta /var/log/letta /tmp

# Environment
Environment=LETTA_CONFIG_PATH=/etc/letta/config.yaml
Environment=PYTHONPATH=/opt/letta

[Install]
WantedBy=multi-user.target
